# infrastructure/docker/services/cron/crontab
0 9 * * * su app -c "/home/app/application/bin/console qotd:run --env=prod" > /proc/1/fd/1 2>&1
