# infrastructure/docker/services/php/cron/crontab
# 0 9 * * * /var/www/bin/console qotd:run > /proc/1/fd/1 2>&1
* * * * * date > /var/log/cron.log 2>&1
