# infrastructure/docker/services/php/cron/crontab
0 9 * * * /home/app/application/bin/console qotd:run > /proc/1/fd/1 2>&1
