# infrastructure/docker/services/php/cron/crontab
* * * * * su app -c "/usr/local/bin/php /var/www/bin/console qotd:run" > /proc/1/fd/1 2>&1
