# infrastructure/docker/services/php/cron/crontab
0 9 * * * /var/www/bin/console qotd:run > /tmp/cron-stdout 2>&1
